# profiles/ubuntu-desktop/root/etc/systemd/system/getty@tty1.service.d/noclear.conf
# https://wiki.archlinux.org/index.php/Getty#Have_boot_messages_stay_on_tty1
[Service]
TTYVTDisallocate=no

